# doc/operation/certbot.timer
[Unit]
Description=Certbot renewal

[Timer]
OnCalendar=Mon *-*-* 00:00:00

[Install]
WantedBy=multi-user.target
